# deploy/data/home_automation.conf
description "the daemon for home automation"
start on starting lighttpd and started postgresql
stop on stopped lighttpd or stopping postgresql
setuid home_automation
setgid home_automation
kill signal INT
exec /home/home_automation/server/lib/sender.py


